# deploy/systemd/link42-agent.service
[Unit]
Description=Link42 Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/link42
Environment=LINK42_SERVER_URL=http://127.0.0.1:8000
Environment=LINK42_NODE_ID=1
Environment=LINK42_AGENT_TOKEN=replace-me
Environment=LINK42_WIREGUARD_DIR=/etc/wireguard
ExecStart=/opt/link42/.venv/bin/python -m link42_agent.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

